# root/etc/uci-defaults/99-l2-config
#!/bin/sh
touch /etc/config/isolation
uci add isolation global
uci set isolation.@global[0].enabled=0
uci add_list isolation.@global[0].networks=''
uci commit isolation
return 0
