# root/etc/uci-defaults/99-l2-config
#!/bin/sh

touch /etc/config/isolation
uci set isolation.enabled=enabled
uci set isolation.networks=networks
uci commit isolation


return 0
